# roles/common/templates/enferno.conf
description "Service to run Gunicorn"
start on runlevel [2345]
stop on runlevel [06]
respawn
respawn limit 10 5
exec /home/{{user_name}}/bin/enferno.sh

